# script/deluser.cron
#!/bin/bash
nas=
dir=
userlist=

while read username; do
	umount /home/$username
	/opt/rocks/bin/rocks run host compute "umount /home/$username"
	/usr/sbin/userdel $username
	sed -i "/$username/ D" /etc/auto.home
	ssh $nas rm -rf $dir/$username
	rm -f /var/spool/mail/$username
done < $userlist
if [[ -s $userlist ]] ; then
	/opt/rocks/bin/rocks sync users > /tmp/cronmylogdel
fi
cat /dev/null > $userlist


